# crontab/park2.php~
#!/usr/bin/php
<?php
echo 'park \n';
